# soulfra.service
[Unit]
Description=Soulfra Ghost Writing Platform
Documentation=https://github.com/calriven/soulfra
After=network.target

[Service]
Type=simple

# User to run as (change to your user)
User=youruser
Group=youruser

# Working directory (change to your installation path)
WorkingDirectory=/var/www/soulfra-simple

# Environment variables
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="FLASK_ENV=production"
Environment="PLATFORM_VERSION=1.0.0"

# Load environment from .env file
EnvironmentFile=-/var/www/soulfra-simple/.env

# Run with gunicorn for production (recommended)
# Install: pip install gunicorn
ExecStart=/usr/bin/gunicorn \
    --workers 4 \
    --bind 0.0.0.0:5001 \
    --timeout 120 \
    --access-logfile /var/log/soulfra/access.log \
    --error-logfile /var/log/soulfra/error.log \
    app:app

# Or run with Flask development server (not recommended for production)
# ExecStart=/usr/bin/python3 /var/www/soulfra-simple/app.py

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/soulfra-simple

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soulfra

[Install]
WantedBy=multi-user.target
